# production/systemd/otel-loki.service
[Unit]
Description=Observability - Loki (Log Backend)
After=docker.service network-online.target
Requires=docker.service

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

ExecStartPre=-/usr/bin/docker rm -f otel-loki
ExecStart=/usr/bin/docker run --rm \
  --name otel-loki \
  --network observability-net \
  -p 127.0.0.1:3100:3100 \
  -v /opt/observability/configs/loki.yaml:/etc/loki/config.yaml:ro \
  -v /opt/observability/data/loki:/loki \
  grafana/loki:3.4.2 \
  -config.file=/etc/loki/config.yaml
ExecStop=/usr/bin/docker stop otel-loki

[Install]
WantedBy=multi-user.target
