# production/systemd/otel-jaeger.service
[Unit]
Description=Observability - Jaeger (Trace Backend)
After=docker.service network-online.target
Requires=docker.service

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

ExecStartPre=-/usr/bin/docker rm -f otel-jaeger
ExecStart=/usr/bin/docker run --rm \
  --name otel-jaeger \
  --network observability-net \
  -p 127.0.0.1:16686:16686 \
  -v /opt/observability/data/jaeger:/badger \
  -e SPAN_STORAGE_TYPE=badger \
  -e BADGER_EPHEMERAL=false \
  -e BADGER_DIRECTORY_VALUE=/badger/data \
  -e BADGER_DIRECTORY_KEY=/badger/key \
  -e COLLECTOR_OTLP_ENABLED=true \
  jaegertracing/all-in-one:1.57
ExecStop=/usr/bin/docker stop otel-jaeger

[Install]
WantedBy=multi-user.target
