# production/systemd/otel-prometheus.service
[Unit]
Description=Observability - Prometheus (Metrics Backend)
After=docker.service network-online.target
Requires=docker.service

[Service]
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30

ExecStartPre=-/usr/bin/docker rm -f otel-prometheus
ExecStart=/usr/bin/docker run --rm \
  --name otel-prometheus \
  --network observability-net \
  -p 127.0.0.1:9090:9090 \
  -v /opt/observability/configs/prometheus.yaml:/etc/prometheus/prometheus.yaml:ro \
  -v /opt/observability/configs/prometheus-web.yaml:/etc/prometheus/web.yaml:ro \
  -v /opt/observability/data/prometheus:/prometheus \
  prom/prometheus:v2.51.2 \
  --config.file=/etc/prometheus/prometheus.yaml \
  --web.config.file=/etc/prometheus/web.yaml \
  --storage.tsdb.path=/prometheus \
  --storage.tsdb.retention.time=30d \
  --storage.tsdb.retention.size=10GB \
  --web.enable-lifecycle
ExecStop=/usr/bin/docker stop otel-prometheus

[Install]
WantedBy=multi-user.target
